# roles/certsync/templates/certsync.service
[Unit]
Description=Synchronize TLS certificates from {{ certsync_host }}

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/certsync
